# debian/manifold.postinst
#!/bin/bash
# if this requires a service to be running, add something like this
# update-rc.d manifold defaults
easy_install python-daemon cfgparse
[ -d /var/myslice ] || mkdir /var/myslice

